# roles/ansible_pull/templates/ansible-pull.timer.j2
# =============================================================================
# IAC-Ansible: ansible-pull systemd timer
# Managed by Ansible - DO NOT EDIT MANUALLY
# =============================================================================

[Unit]
Description=IAC-Ansible Pull Timer - Check for configuration updates

[Timer]
OnCalendar={{ ansible_pull_schedule }}
RandomizedDelaySec={{ ansible_pull_random_delay_sec }}
Persistent=true
AccuracySec=60

[Install]
WantedBy=timers.target
